# {{cookiecutter.app_name}}/playbooks/templates/app.ini
[program:]
command=/usr/local/bin/uwsgi
  -- 0.0.0.0:
  --logto /log/application.log
  --venv /venv
  --pythonpath /src
  --wsgi-file /src/wsgi.py
  --callable application
  --max-requests 1000
  --master
  --processes 1
  --chmod
directory=/src
autostart=true
autorestart=true