# {{cookiecutter.app_name}}/playbooks/templates/supervisor/{{cookiecutter.app_name}}.ini
{%- raw -%}
[program:{{ app_name }}]
command={{ app_path }}/venv/bin/gunicorn
    --name {{ app_name }}.gunicorn
    --bind unix:{{ app_instance_path }}/{{ app_name }}.socket
    --workers 1
    --threads 1
    --log-level info
    --env {{ app_name|upper }}_INSTANCE_PATH={{ app_instance_path }}
    {% if app_debug -%} --reload {% endif -%}
    manage:application
directory={{ app_path }}/src
{%- endraw -%}
