# {{cookiecutter.app_name}}/playbooks/roles/gateone/templates/gateone.ini
{%- raw -%}
[program:gateone]
command={{ gateone_path }}/venv/bin/gateone
{%- endraw -%}
